# infra/linux/duallink-receiver.service
[Unit]
Description=DualLink Receiver — macOS display extension for Linux
Documentation=https://github.com/cassianoheitor1993/DualLink
After=network-online.target graphical-session.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/duallink-receiver
Restart=on-failure
RestartSec=3s

# Allow GStreamer to open a display window
Environment=DISPLAY=:0
Environment=XDG_RUNTIME_DIR=/run/user/1000
Environment=WAYLAND_DISPLAY=wayland-0

# Log via systemd journal
StandardOutput=journal
StandardError=journal
SyslogIdentifier=duallink-receiver

# Reasonable resource limits
LimitNOFILE=65536
Nice=-5

[Install]
WantedBy=graphical-session.target
